# crontab
* * * * * ubuntu /bin/bash /home/ubuntu/apps/rockymountain/repository/run_on_server.sh fab $DEPLOYMENT_TARGET cron_jobs.test >/var/log/rockymountain/crontab.log 2>&1
